# root/etc/services.d/dockupdater/run
#!/usr/bin/with-contenv bash

cd /app/dockupdater || exit

exec \
	s6-setuidgid root python3 entrypoint
