# root/etc/services.d/dockupdate/run
#!/usr/bin/with-contenv bash

cd /app/dockupdate || exit

exec \
	s6-setuidgid root python3 entrypoint
